# overlay/etc/init.d/20_symlink
#!/sbin/openrc-run

command="ln -s /run /var/run"
