# overlay/etc/init.d/20_symlink
#!/sbin/openrc-run

name="20_symlink"
description="Create symlink /var/run -> /run"

depend() {
    after 00_link_mount
}

start() {
    ln -s /run /var/run
}
